# trinitycore-tools3.3.5.files/B_DEBIAN/etc/profile.d/trinitycore.sh
#!/bin/bash
# Set PATH variable for TrinityCore

PATH=/opt/trinitycore/bin:$PATH
